# systemd/battycoda-celery.service
[Unit]
Description=BattyCoda Celery Worker
After=network.target redis.service

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/battycoda
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=0
ExecStart=/usr/local/bin/celery -A celery_app.celery worker --loglevel=info
Restart=always
RestartSec=5

# Limit resources
MemoryLimit=1G
CPUQuota=30%

[Install]
WantedBy=multi-user.target