# systemd/battycoda-r-direct.service
[Unit]
Description=BattyCoda R Direct Server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/battycoda
ExecStart=/usr/bin/Rscript r_server_direct.R
Restart=always
RestartSec=5

# Limit resources
MemoryLimit=1G
CPUQuota=30%

[Install]
WantedBy=multi-user.target